# releng/airootfs/etc/skel/.config/autostart/start.desktop
[Desktop Entry]
Version=1.0
Name=start
Exec=core-start
Type=Application
Terminal=false
Categories=Utility
